# etc/supervisord.d/crond.ini
[program:crond]
command=/usr/sbin/crond -n
numprocs=1
autostart=true
autorestart=true
