# etc/supervisord.d/dovecot.ini
[program:dovecot]
command=/usr/sbin/dovecot -F
numprocs=1
autostart=true
autorestart=true
