# etc/supervisord.d/httpd.ini
[program:httpd]
command=/sbin/httpd -DFOREGROUND
numprocs=1
autostart=true
autorestart=true
