# etc/supervisord.d/mysqld.ini
[program:mysqld]
command=/usr/sbin/mysqld --pid-file=/var/run/mysqld/mysqld.pid --user=root
numprocs=1
autostart=true
autorestart=true
